# agents/granola-google-docs-importer/config.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.petesbrain.granola-google-docs-importer</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>set -a; source /Users/administrator/Documents/PetesBrain.nosync/.env; set +a; /Users/administrator/Documents/PetesBrain.nosync/shared/email-sync/.venv/bin/python3 /Users/administrator/Documents/PetesBrain.nosync/agents/granola-google-docs-importer/granola-google-docs-importer.py --days 7</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict/>
    <key>StartInterval</key>
    <integer>3600</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/administrator/.petesbrain-granola-google-docs.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/administrator/.petesbrain-granola-google-docs.error.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/administrator/Documents/PetesBrain.nosync</string>
    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
