# infrastructure/rollback-snapshots/20251211_162054/plist/com.petesbrain.populate-spreadsheets.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>GOOGLE_APPLICATION_CREDENTIALS</key>
		<string>/Users/administrator/Documents/PetesBrain.nosync/infrastructure/mcp-servers/google-sheets-mcp-server/credentials.json</string>
		<key>PATH</key>
		<string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
	</dict>
	<key>Label</key>
	<string>com.petesbrain.populate-spreadsheets</string>
	<key>ProcessType</key>
	<string>Background</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/local/bin/python3</string>
		<string>/Users/administrator/Documents/PetesBrain.nosync/tools/product-impact-analyzer/populate_all_spreadsheets.py</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/Users/administrator/.petesbrain-populate-spreadsheets.log</string>
	<key>StandardOutPath</key>
	<string>/Users/administrator/.petesbrain-populate-spreadsheets.log</string>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>8</integer>
		<key>Minute</key>
		<integer>15</integer>
	</dict>
	<key>WorkingDirectory</key>
	<string>/Users/administrator/Documents/PetesBrain.nosync/tools/product-impact-analyzer</string>
</dict>
</plist>
