# infrastructure/rollback-snapshots/20251211_162054/plist/com.petesbrain.email-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PETESBRAIN_ROOT</key>
		<string>/Users/administrator/Documents/PetesBrain.nosync</string>
		<key>GOOGLE_APPLICATION_CREDENTIALS</key>
		<string>/Users/administrator/Documents/PetesBrain.nosync/shared/email-sync/credentials.json</string>
	</dict>
	<key>KeepAlive</key>
	<false/>
	<key>Label</key>
	<string>com.petesbrain.email-sync</string>
	<key>ProcessType</key>
	<string>Background</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/local/bin/python3</string>
		<string>/Users/administrator/Documents/PetesBrain.nosync/shared/email-sync/email_sync_workflow.py</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/Users/administrator/.petesbrain-email-sync-error.log</string>
	<key>StandardOutPath</key>
	<string>/Users/administrator/.petesbrain-email-sync.log</string>
	<key>StartInterval</key>
	<integer>21600</integer>
	<key>WorkingDirectory</key>
	<string>/Users/administrator/Documents/PetesBrain.nosync</string>
</dict>
</plist>
