# infrastructure/rollback-snapshots/20251213_080004/plist/com.petesbrain.weekly-label-reports.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.petesbrain.weekly-label-reports</string>
	<key>ProcessType</key>
	<string>Background</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>-c</string>
		<string>source $HOME/.bashrc 2&gt;/dev/null || source $HOME/.zshrc 2&gt;/dev/null; /usr/local/bin/python3 /Users/administrator/Documents/PetesBrain.nosync/tools/product-impact-analyzer/label_validation_report.py --send-email</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/Users/administrator/.petesbrain-weekly-label-reports-error.log</string>
	<key>StandardOutPath</key>
	<string>/Users/administrator/.petesbrain-weekly-label-reports.log</string>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>9</integer>
		<key>Minute</key>
		<integer>0</integer>
		<key>Weekday</key>
		<integer>1</integer>
	</dict>
	<key>WorkingDirectory</key>
	<string>/Users/administrator/Documents/PetesBrain.nosync/tools/product-impact-analyzer</string>
</dict>
</plist>
