# agents/tasks-backup/config.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.petesbrain.tasks-backup</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/python3</string>
        <string>/Users/administrator/Documents/PetesBrain/agents/tasks-backup/tasks-backup.py</string>
    </array>

    <key>StartInterval</key>
    <integer>10800</integer>

    <key>StandardOutPath</key>
    <string>/Users/administrator/.petesbrain-tasks-backup.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/administrator/.petesbrain-tasks-backup-error.log</string>

    <key>RunAtLoad</key>
    <true/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>PYTHONPATH</key>
        <string>/Users/administrator/Documents/PetesBrain</string>
    </dict>
</dict>
</plist>
