# tools/product-impact-analyzer/launchagents/com.petesbrain.product-tracking.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.petesbrain.product-tracking</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>cd /Users/administrator/Documents/PetesBrain/tools/product-impact-analyzer &amp;&amp; GOOGLE_APPLICATION_CREDENTIALS=/Users/administrator/Documents/PetesBrain/infrastructure/mcp-servers/google-sheets-mcp-server/credentials.json .venv/bin/python3 product_feed_tracker.py &amp;&amp; sleep 60 &amp;&amp; .venv/bin/python3 product_change_detector.py</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/administrator/.petesbrain-product-tracking.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/administrator/.petesbrain-product-tracking.log</string>
</dict>
</plist>
