# tools/product-impact-analyzer/launchagents/com.petesbrain.baseline-calculator.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.petesbrain.baseline-calculator</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>cd /Users/administrator/Documents/PetesBrain/tools/product-impact-analyzer &amp;&amp; GOOGLE_APPLICATION_CREDENTIALS=/Users/administrator/Documents/PetesBrain/infrastructure/mcp-servers/google-sheets-mcp-server/credentials.json .venv/bin/python3 product_baseline_calculator.py</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/administrator/.petesbrain-baseline-calculator.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/administrator/.petesbrain-baseline-calculator.log</string>
</dict>
</plist>
