# shared/scripts/com.petesbrain.budget-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.petesbrain.budget-monitor</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/administrator/Documents/PetesBrain/shared/email-sync/.venv/bin/python3</string>
        <string>/Users/administrator/Documents/PetesBrain/shared/scripts/daily-budget-monitor.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/administrator/Documents/PetesBrain</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>GOOGLE_APPLICATION_CREDENTIALS</key>
        <string>/Users/administrator/Documents/PetesBrain/shared/email-sync/credentials.json</string>
    </dict>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/administrator/.petesbrain-budget-monitor.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/administrator/.petesbrain-budget-monitor.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
