# infrastructure/rollback-snapshots/20251211_143600/plist/com.petesbrain.task-notes-server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PYTHONDONTWRITEBYTECODE</key>
		<string>1</string>
		<key>PYTHONUNBUFFERED</key>
		<string>1</string>
		<key>PYTHONNOUSERSITE</key>
		<string>1</string>
	</dict>
	<key>Label</key>
	<string>com.petesbrain.task-notes-server</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/local/bin/python3</string>
		<string>/Users/administrator/Documents/PetesBrain.nosync/shared/scripts/save-task-notes.py</string>
	</array>
	<key>KeepAlive</key>
	<true/>
	<key>StandardOutPath</key>
	<string>/Users/administrator/.petesbrain-task-notes-server.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/administrator/.petesbrain-task-notes-server-error.log</string>
	<key>StartInterval</key>
	<integer>0</integer>
</dict>
</plist>
