# infrastructure/rollback-snapshots/20251211_161906/plist/com.petesbrain.weekly-impact-report.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.petesbrain.weekly-impact-report</string>
	<key>ProcessType</key>
	<string>Background</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>-c</string>
		<string>cd /Users/administrator/Documents/PetesBrain.nosync/tools/product-impact-analyzer &amp;&amp; GOOGLE_APPLICATION_CREDENTIALS=/Users/administrator/Documents/PetesBrain.nosync/infrastructure/mcp-servers/google-sheets-mcp-server/credentials.json /usr/local/bin/python3 weekly_impact_report.py</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/Users/administrator/.petesbrain-weekly-impact-report.log</string>
	<key>StandardOutPath</key>
	<string>/Users/administrator/.petesbrain-weekly-impact-report.log</string>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>9</integer>
		<key>Minute</key>
		<integer>15</integer>
		<key>Weekday</key>
		<integer>1</integer>
	</dict>
</dict>
</plist>
