# infrastructure/rollback-snapshots/20251211_161906/plist/com.petesbrain.critical-tasks-backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.petesbrain.critical-tasks-backup</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>
set -e
TIMESTAMP=$(date +%Y%m%d-%H%M%S)
ICLOUD_DIR="$HOME/Library/Mobile Documents/com~apple~CloudDocs/PetesBrain-Backups/critical-tasks"
BACKUP_NAME="tasks-backup-${TIMESTAMP}.tar.gz"
mkdir -p "$ICLOUD_DIR"
cd /Users/administrator/Documents/PetesBrain.nosync
tar -czf "$ICLOUD_DIR/$BACKUP_NAME" clients/*/tasks.json clients/*/tasks-completed.md roksys/spreadsheets/rok-experiments-client-notes.csv 2>/dev/null || true
ls -1t "$ICLOUD_DIR"/tasks-backup-*.tar.gz 2>/dev/null | tail -n +21 | xargs rm -f 2>/dev/null || true
echo "[$(date)] Critical tasks backed up: $BACKUP_NAME"
        </string>
    </array>

    <key>StartInterval</key>
    <integer>21600</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/administrator/.petesbrain-tasks-backup.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/administrator/.petesbrain-tasks-backup-error.log</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
