# agents/whatsapp-processor/config.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.petesbrain.whatsapp-processor</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/python3</string>
        <string>/Users/administrator/Documents/PetesBrain/agents/whatsapp-processor/whatsapp-processor.py</string>
        <string>--days</string>
        <string>7</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/administrator/Documents/PetesBrain</string>

    <key>StartInterval</key>
    <integer>1800</integer>

    <key>StandardOutPath</key>
    <string>/Users/administrator/.petesbrain-whatsapp.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/administrator/.petesbrain-whatsapp-error.log</string>

    <key>RunAtLoad</key>
    <true/>

<key>ProcessType</key>
<string>Background</string>
</dict>
</plist>
