# infrastructure/rollback-snapshots/20251211_150404/plist/com.petesbrain.daily-anomaly-alerts.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>GOOGLE_ADS_CONFIGURATION_FILE_PATH</key>
		<string>/Users/administrator/google-ads.yaml</string>
	</dict>
	<key>Label</key>
	<string>com.petesbrain.daily-anomaly-alerts</string>
	<key>ProcessType</key>
	<string>Background</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/local/bin/python3</string>
		<string>/Users/administrator/Documents/PetesBrain.nosync/agents/performance-monitoring/daily-anomaly-detector.py</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/Users/administrator/.petesbrain-daily-anomaly.log</string>
	<key>StandardOutPath</key>
	<string>/Users/administrator/.petesbrain-daily-anomaly.log</string>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>6</integer>
		<key>Minute</key>
		<integer>30</integer>
	</dict>
	<key>WorkingDirectory</key>
	<string>/Users/administrator/Documents/PetesBrain.nosync</string>
</dict>
</plist>
