# infrastructure/rollback-snapshots/20251212_080000/plist/com.petesbrain.task-manager-hourly-regenerate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.petesbrain.task-manager-hourly-regenerate</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/local/bin/python3</string>
		<string>/Users/administrator/Documents/PetesBrain.nosync/shared/scripts/generate-task-manager.py</string>
	</array>
	<key>StartInterval</key>
	<integer>3600</integer>
	<key>StandardOutPath</key>
	<string>/Users/administrator/.petesbrain-task-manager-regenerate.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/administrator/.petesbrain-task-manager-regenerate-error.log</string>
</dict>
</plist>
