# infrastructure/rollback-snapshots/20251212_080000/plist/com.petesbrain.agent-loader.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.petesbrain.agent-loader</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>
LOG_FILE="$HOME/.petesbrain-agent-loader.log"
LAUNCH_AGENTS_DIR="$HOME/Library/LaunchAgents"
echo "============================================================" >> "$LOG_FILE"
echo "$(date '+%Y-%m-%d %H:%M:%S') - Agent Loader Starting (Staggered)" >> "$LOG_FILE"
sleep 10
LOADED_COUNT=0
ALREADY_LOADED=0
for plist in "$LAUNCH_AGENTS_DIR"/com.petesbrain.*.plist; do
    if [ -f "$plist" ]; then
        LABEL=$(basename "$plist" .plist)

        # Skip the agent-loader itself to prevent recursion
        if [[ "$LABEL" == "com.petesbrain.agent-loader" ]]; then
            continue
        fi

        if launchctl list "$LABEL" &amp;>/dev/null; then
            ((ALREADY_LOADED++))
            continue
        fi

        if launchctl load "$plist" 2>/dev/null; then
            echo "Loaded: $LABEL" >> "$LOG_FILE"
            ((LOADED_COUNT++))
            # Add 2-second delay between each agent to prevent resource contention
            sleep 2
        fi
    fi
done
echo "Already loaded: $ALREADY_LOADED, Newly loaded: $LOADED_COUNT" >> "$LOG_FILE"
echo "$(date '+%Y-%m-%d %H:%M:%S') - Agent Loader Complete" >> "$LOG_FILE"
        </string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>LimitLoadToSessionType</key>
    <string>Aqua</string>

    <key>StandardOutPath</key>
    <string>/Users/administrator/.petesbrain-agent-loader.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/administrator/.petesbrain-agent-loader-error.log</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
