# infrastructure/rollback-snapshots/20251211_161833/plist/com.petesbrain.granola-weekly-summary.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.petesbrain.granola-weekly-summary</string>
	<key>ProcessType</key>
	<string>Background</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Users/administrator/Documents/PetesBrain.nosync/tools/granola-importer/venv/bin/python3</string>
		<string>/Users/administrator/Documents/PetesBrain.nosync/tools/granola-importer/send_weekly_summary.py</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/Users/administrator/.petesbrain-granola-weekly-email.log</string>
	<key>StandardOutPath</key>
	<string>/Users/administrator/.petesbrain-granola-weekly-email.log</string>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>9</integer>
		<key>Minute</key>
		<integer>0</integer>
		<key>Weekday</key>
		<integer>1</integer>
	</dict>
	<key>WorkingDirectory</key>
	<string>/Users/administrator/Documents/PetesBrain.nosync/tools/granola-importer</string>
</dict>
</plist>
