# tools/product-impact-analyzer/com.petesbrain.product-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.petesbrain.product-monitor</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/administrator/Documents/PetesBrain/tools/product-impact-analyzer/.venv/bin/python3</string>
        <string>/Users/administrator/Documents/PetesBrain/tools/product-impact-analyzer/monitor.py</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>10</integer>  <!-- 10 AM daily -->
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>EnvironmentVariables</key>
    <dict>
        <key>ANTHROPIC_API_KEY</key>
        <string>WILL_BE_REPLACED_BY_SETUP_SCRIPT</string>
        <key>GOOGLE_APPLICATION_CREDENTIALS</key>
        <string>/Users/administrator/Documents/PetesBrain/infrastructure/mcp-servers/google-sheets-mcp-server/credentials.json</string>
        <key>GMAIL_APP_PASSWORD</key>
        <string>WILL_BE_REPLACED_BY_SETUP_SCRIPT</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/administrator/Documents/PetesBrain/tools/product-impact-analyzer</string>

    <key>StandardOutPath</key>
    <string>/Users/administrator/.petesbrain-product-monitor.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/administrator/.petesbrain-product-monitor-error.log</string>

    <key>RunAtLoad</key>
    <false/>  <!-- Don't run on load, only on schedule -->

    <key>Nice</key>
    <integer>1</integer>  <!-- Lower priority -->
</dict>
</plist>
