# infrastructure/rollback-snapshots/20251215_080000/plist/com.petesbrain.sync-todos-to-google-tasks.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>GOOGLE_APPLICATION_CREDENTIALS</key>
		<string>/Users/administrator/Documents/PetesBrain.nosync/infrastructure/mcp-servers/google-tasks-mcp-server/credentials.json</string>
	</dict>
	<key>Label</key>
	<string>com.petesbrain.sync-todos-to-google-tasks</string>
	<key>ProcessType</key>
	<string>Background</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/local/bin/python3</string>
		<string>/Users/administrator/Documents/PetesBrain.nosync/agents/sync-todos-to-google-tasks/sync-todos-to-google-tasks.py</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/Users/administrator/.petesbrain-sync-todos.log</string>
	<key>StandardOutPath</key>
	<string>/Users/administrator/.petesbrain-sync-todos.log</string>
	<key>StartInterval</key>
	<integer>3600</integer>
	<key>WorkingDirectory</key>
	<string>/Users/administrator/Documents/PetesBrain</string>
</dict>
</plist>
