# infrastructure/rollback-snapshots/20251215_080000/plist/com.petesbrain.tree2mydoor-context-upload.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.petesbrain.tree2mydoor-context-upload</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>
CONTEXT_FILE="/Users/administrator/Documents/PetesBrain.nosync/clients/tree2mydoor/CONTEXT.md"
SYNC_DIR="$HOME/Library/CloudStorage/GoogleDrive-petere@roksys.co.uk/My Drive/PetesBrain-Context"
LOG_FILE="$HOME/.petesbrain-tree2mydoor-context-upload.log"
DATE=$(date +"%Y-%m-%d")
TIMESTAMP=$(date +"%Y-%m-%d %H:%M:%S")

echo "[$TIMESTAMP] Tree2MyDoor CONTEXT.md Daily Sync" >> "$LOG_FILE"
mkdir -p "$SYNC_DIR"

if [ ! -f "$CONTEXT_FILE" ]; then
    echo "[$TIMESTAMP] ERROR: CONTEXT.md not found" >> "$LOG_FILE"
    exit 1
fi

DEST_FILE="$SYNC_DIR/Tree2MyDoor-CONTEXT-$DATE.md"
cp "$CONTEXT_FILE" "$DEST_FILE"
ln -sf "$DEST_FILE" "$SYNC_DIR/Tree2MyDoor-CONTEXT-LATEST.md"

echo "[$TIMESTAMP] Synced to: $DEST_FILE" >> "$LOG_FILE"
        </string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/administrator/.petesbrain-tree2mydoor-context-upload.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/administrator/.petesbrain-tree2mydoor-context-upload-error.log</string>

    <key>ProcessType</key>
    <string>Background</string>

    <key>WorkingDirectory</key>
    <string>/Users/administrator/Documents/PetesBrain.nosync</string>
</dict>
</plist>
